# config/systemd/user/ArchUpdates.service
[Unit]
Description=Check for available updates & sync polybar or eww bar

[Service]
Type=oneshot
ExecStart=%h/.config/bspwm/src/Updates --sync-polybar

[Install]
WantedBy=default.target

